# puppet-manifests/src/modules/platform/templates/phc2sys.service.erb
[Unit]
Description=Synchronize system clock or PTP hardware clock (PHC)
After=ntpdate.service

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/phc2sys
ExecStart=/usr/sbin/phc2sys $OPTIONS
ExecStartPost=/bin/bash -c 'echo $MAINPID > /var/run/phc2sys.pid'
ExecStopPost=/bin/rm -f /var/run/phc2sys.pid

[Install]
WantedBy=multi-user.target
