# puppet-manifests/src/modules/platform/templates/ptp4l.service.erb
[Unit]
Description=Precision Time Protocol (PTP) service

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/ptp4l
ExecStart=/usr/sbin/ptp4l $OPTIONS
ExecStartPost=/bin/bash -c 'echo $MAINPID > /var/run/ptp4l.pid'
ExecStopPost=/bin/rm -f /var/run/ptp4l.pid

[Install]
WantedBy=multi-user.target
