# sysinv/sysinv-fpga-agent/sysinv-conf-watcher.path
[Path]
PathChanged=/etc/sysinv/sysinv.conf

[Install]
WantedBy=multi-user.target
